# infra/quadlet-sync/quadlet-sync.timer
[Unit]
Description=Sync quadlet configs from git

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
